# platform/host/certificates/templates/wapps-acme-renew.service.j2
[Unit]
Description=Wapps ACME certificate renewal (acme.sh)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/wapps/acme.env
ExecStart=/root/.acme.sh/acme.sh --cron --home /root/.acme.sh

